# systemd/mce.service
[Unit]
Description=Mode Control Entity (MCE)
Requires=dsme.service
Requires=dbus.service
After=dsme.service
# We want following but for some reason it causes shutdown problem
# disable until we find out why
# Before=systemd-user-sessions.service

[Service]
Type=notify
ExecStart=/usr/sbin/mce --systemd
Restart=always

[Install]
WantedBy=multi-user.target
