# systemd/mce.service
[Unit]
Description=Mode Control Entity (MCE)
Requires=dsme.service
Requires=dbus.service
After=dsme.service
Before=start-user-session@USER.service start-user-session@ACT_DEAD.service

[Service]
Type=notify
ExecStart=/usr/sbin/mce --systemd
Restart=always

[Install]
WantedBy=multi-user.target
