# systemd/mce.service
[Unit]
Description=Mode Control Entity (MCE)
Requires=dsme.service
Requires=dbus.service
After=dsme.service
Before=systemd-user-sessions.service

[Service]
ExecStart=/usr/sbin/mce
Restart=always

[Install]
WantedBy=multi-user.target
